# home/.config/autostart/anacron.sh.desktop
[Desktop Entry]
Type=Application
Name=anacron.sh
Exec=bash -c '${HOME}/bin/anacron.sh'
